# deploy/supervisor/monitly-queues.conf
[group:monitly]
programs=monitly-checks-priority,monitly-checks-standard,monitly-incidents,monitly-notifications,monitly-webhooks-in,monitly-webhooks-out,monitly-sla,monitly-maintenance,monitly-exports

[program:monitly-checks-priority]
process_name=%(program_name)s_%(process_num)02d
command=php /workspace/monitly-app/artisan queue:work redis --queue=checks_priority --sleep=1 --tries=3 --timeout=30
numprocs=4
autostart=true
autorestart=true
stopwaitsecs=360
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/monitly/checks-priority.log
stopasgroup=true
killasgroup=true
environment=APP_ENV="production",QUEUE_CONNECTION="redis"

[program:monitly-checks-standard]
process_name=%(program_name)s_%(process_num)02d
command=php /workspace/monitly-app/artisan queue:work redis --queue=checks_standard --sleep=1 --tries=3 --timeout=30
numprocs=2
autostart=true
autorestart=true
stopwaitsecs=360
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/monitly/checks-standard.log
stopasgroup=true
killasgroup=true
environment=APP_ENV="production",QUEUE_CONNECTION="redis"

[program:monitly-incidents]
process_name=%(program_name)s_%(process_num)02d
command=php /workspace/monitly-app/artisan queue:work redis --queue=incidents --sleep=1 --tries=3 --timeout=30
numprocs=1
autostart=true
autorestart=true
stopwaitsecs=360
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/monitly/incidents.log
stopasgroup=true
killasgroup=true
environment=APP_ENV="production",QUEUE_CONNECTION="redis"

[program:monitly-notifications]
process_name=%(program_name)s_%(process_num)02d
command=php /workspace/monitly-app/artisan queue:work redis --queue=notifications --sleep=1 --tries=5 --timeout=60
numprocs=2
autostart=true
autorestart=true
stopwaitsecs=360
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/monitly/notifications.log
stopasgroup=true
killasgroup=true
environment=APP_ENV="production",QUEUE_CONNECTION="redis"

[program:monitly-webhooks-in]
process_name=%(program_name)s_%(process_num)02d
command=php /workspace/monitly-app/artisan queue:work redis --queue=webhooks_in --sleep=1 --tries=5 --timeout=30
numprocs=1
autostart=true
autorestart=true
stopwaitsecs=360
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/monitly/webhooks-in.log
stopasgroup=true
killasgroup=true
environment=APP_ENV="production",QUEUE_CONNECTION="redis"

[program:monitly-webhooks-out]
process_name=%(program_name)s_%(process_num)02d
command=php /workspace/monitly-app/artisan queue:work redis --queue=webhooks_out --sleep=1 --tries=8 --timeout=30
numprocs=1
autostart=true
autorestart=true
stopwaitsecs=360
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/monitly/webhooks-out.log
stopasgroup=true
killasgroup=true
environment=APP_ENV="production",QUEUE_CONNECTION="redis"

[program:monitly-sla]
process_name=%(program_name)s_%(process_num)02d
command=php /workspace/monitly-app/artisan queue:work redis --queue=sla --sleep=5 --tries=3 --timeout=120
numprocs=1
autostart=true
autorestart=true
stopwaitsecs=360
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/monitly/sla.log
stopasgroup=true
killasgroup=true
environment=APP_ENV="production",QUEUE_CONNECTION="redis"

[program:monitly-maintenance]
process_name=%(program_name)s_%(process_num)02d
command=php /workspace/monitly-app/artisan queue:work redis --queue=maintenance --sleep=5 --tries=3 --timeout=120
numprocs=1
autostart=true
autorestart=true
stopwaitsecs=360
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/monitly/maintenance.log
stopasgroup=true
killasgroup=true
environment=APP_ENV="production",QUEUE_CONNECTION="redis"

[program:monitly-exports]
process_name=%(program_name)s_%(process_num)02d
command=php /workspace/monitly-app/artisan queue:work redis --queue=exports --sleep=5 --tries=3 --timeout=120
numprocs=1
autostart=true
autorestart=true
stopwaitsecs=360
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/monitly/exports.log
stopasgroup=true
killasgroup=true
environment=APP_ENV="production",QUEUE_CONNECTION="redis"
